# arkdep-build.d/biont/overlay/post_install/etc/systemd/system/arkane-postinstall.service
[Unit]
Description=Run arkane-postinstall once machine is online
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/arkane-postinstall
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
