# arkdep-build.d/biont/overlay/post_install/etc/systemd/system/systemd-gig@.socket
[Unit]
PartOf=server-online.target
[Socket]
ListenStream=%i
[Install]
WantedBy=sockets.target
WantedBy=server-online.target
